# ForBarchart/crontab
10 5-22 * * * cd /forbarchart_cronjob && /usr/local/bin/python forbarchart.py >> cron.log 2>&1
