# Evaluate_report_all/crontab
10 1,9,12,18 * * * cd /report_all_cronjob && /usr/local/bin/python evaluate_report_all.py >> cron.log 2>&1
