# ReportForZone/ReportZone11/crontab
56 5,10,16,19 * * * cd /report_zone11_cronjob && /usr/local/bin/python report_zone11.py >> cron.log 2>&1
